# meta-userboard-g2l/recipes-demo/usbcam-http-demo/files/thttpd-demo.service
[Unit]
Description=USBCAM HTTPD DEMO
After=multi-user.target

[Service]
Environment=XDG_RUNTIME_DIR=/run/user/0
ExecStart=/home/root/tvm/thttpd-demo.sh
Type=simple

[Install]
WantedBy=multi-user.target
